# infra/optimal-update.timer
[Unit]
Description=Check for OpenClaw updates 6x daily

[Timer]
OnCalendar=*-*-* 02,06,10,14,18,22:00:00
Persistent=true

[Install]
WantedBy=timers.target
