# infra/optimal-iteration.service
[Unit]
Description=Pick up and work on next Supabase task
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=oracle
WorkingDirectory=/home/oracle/optimal-cli
ExecStart=/home/oracle/.npm-global/bin/tsx bin/optimal.ts agent claim --agent oracle
TimeoutStartSec=600
EnvironmentFile=/home/oracle/optimal-cli/.env
Environment=HOME=/home/oracle
Environment=PATH=/home/oracle/.npm-global/bin:/home/oracle/.bun/bin:/usr/local/bin:/usr/bin:/bin
